# .ebextensions/crontab
0 0 * * * /usr/local/bin/node npm run scrape